# services/linting/markdown-manager-linting-consumer.service
# Systemd service configuration for markdown-lint-consumer
# Place in /etc/systemd/system/markdown-manager-lint-consumer.service

[Unit]
Description=Markdown Lint Consumer Service for Markdown Manager
After=network.target docker.service redis.service postgresql.service
Wants=network.target
Requires=docker.service
Wants=redis.service postgresql.service

[Service]
Type=simple
User=dlittle
Group=dlittle
EnvironmentFile=-/etc/markdown-manager.env
ExecStartPre=-/usr/bin/docker stop markdown-manager-linting-consumer
ExecStartPre=-/usr/bin/docker rm markdown-manager-linting-consumer
ExecStart=/usr/bin/docker run --rm --name markdown-manager-linting-consumer \
  -e DATABASE_URL=${DATABASE_URL} \
  -e REDIS_URL=${REDIS_URL} \
  -v /opt/markdown-manager/configs/linting-consumer.config.json:/app/config/consumer.config.json:ro \
  --health-cmd="python -c 'import redis.asyncio; import asyncio; asyncio.run(redis.asyncio.from_url(\"${REDIS_URL}\").ping())'" \
  --health-interval=30s \
  --health-timeout=10s \
  --health-retries=3 \
  --health-start-period=20s \
  littledan9/markdown-manager-event-consumer:latest
ExecStop=/usr/bin/docker stop markdown-manager-linting-consumer
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=markdown-manager-linting-consumer

[Install]
WantedBy=multi-user.target