# backend/markdown-manager-api.service

[Unit]
Description=Markdown Manager API (Docker)
After=network.target docker.service
Requires=docker.service


[Service]
Type=simple
Restart=on-failure
RestartSec=5s
EnvironmentFile=/etc/markdown-manager.env
ExecStart=/usr/bin/docker run --rm --name markdown-manager-api \
  -e DATABASE_URL=${DATABASE_URL} \
  -e EXPORT_SERVICE_URL=${EXPORT_SERVICE_URL} \
  -e MARKDOWN_LINT_SERVICE_URL=${MARKDOWN_LINT_SERVICE_URL} \
  -e SPELL_CHECK_SERVICE_URL=${SPELL_CHECK_SERVICE_URL} \
  -e ENVIRONMENT=${ENVIRONMENT} \
  -e GITHUB_CLIENT_ID=${GITHUB_CLIENT_ID} \
  -e GITHUB_CLIENT_SECRET=${GITHUB_CLIENT_SECRET} \
  -e GITHUB_REDIRECT_URI=${GITHUB_REDIRECT_URI} \
  -e GITHUB_OAUTH_SCOPE=${GITHUB_OAUTH_SCOPE} \
  -e GITHUB_MAX_REPOS_PER_ACCOUNT=${GITHUB_MAX_REPOS_PER_ACCOUNT} \
  -e GITHUB_MIN_UPDATED_DAYS=${GITHUB_MIN_UPDATED_DAYS} \
  -e GITHUB_INCLUDE_FORKS=${GITHUB_INCLUDE_FORKS} \
  -e GITHUB_EXCLUDE_ARCHIVED=${GITHUB_EXCLUDE_ARCHIVED} \
  -e JWT_SECRET_KEY=${JWT_SECRET_KEY} \
  -e MARKDOWN_STORAGE_ROOT=${CONTAINER_STORAGE_ROOT} \
  -v ${HOST_STORAGE_ROOT}:/documents \
  -p 8000:8000 \
  littledan9/markdown-manager:latest
ExecStop=/usr/bin/docker stop markdown-manager-api


[Install]
WantedBy=multi-user.target
