# export-service/markdown-manager-export.service
[Unit]
Description=Markdown Manager Export Service (Docker)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
Restart=on-failure
RestartSec=5s
EnvironmentFile=/etc/markdown-manager.env
ExecStart=/usr/bin/docker run --rm --name markdown-manager-export \
  -e PYTHONUNBUFFERED=1 \
  -p 8001:8001 \
  littledan9/markdown-manager-export:latest
ExecStop=/usr/bin/docker stop markdown-manager-export

[Install]
WantedBy=multi-user.target
